# conf/dingus.upstart.conf
# Place this file in /etc/init/dingus.conf
description "dingus"
author "Ben Reichelt"

start on net-device-up IFACE=eth0

script
	sleep 20
	chdir /home/web/dingus
	exec ./dingus
end script
